# demon2.service
#!/bin/bash
[Unit]
Description=El_señor_script

[Service]
ExecStart=/usr/bin/demon

[Install]
WantedBy=multi-user.target 
